# systemd/suse-migration-reboot.service
[Unit]
Description=Reboot System
After=suse-migration-kernel-load.service suse-migration-update-bootloader.service
ConditionKernelCommandLine=!migration.noreboot

[Service]
Type=oneshot
ExecStart=/usr/bin/suse-migration-reboot

[Install]
WantedBy=multi-user.target
