# systemd/suse-migration-container-emergency.service
[Unit]
Description=Reboot System

[Service]
Type=oneshot
ExecStart=/usr/bin/suse-migration-reboot

[Install]
WantedBy=multi-user.target
