# systemd/suse-migration.service
[Unit]
Description=Run Zypper Migration
After=suse-migration-product-setup.service
Requires=suse-migration-product-setup.service

[Service]
Type=oneshot
ExecStart=/usr/bin/suse-migration

[Install]
WantedBy=multi-user.target
